# roles/homelab-apps/templates/refresh-secrets.service.j2
[Unit]
Description=Refresh Homelab App Secrets from Infisical
Documentation=https://github.com/LuckyPrayer/homelab-apps
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/refresh-secrets.sh
User=root
Group=root

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=refresh-secrets

[Install]
WantedBy=multi-user.target
