# roles/homelab-compose/templates/homelab-stacks.service.j2
[Unit]
Description=Homelab Docker Stacks
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=docker-svc
Group=docker
WorkingDirectory=/opt/stacks
Environment=INFISICAL_ENV={{ hostvars[inventory_hostname].env | default('dev') }}
Environment=INFISICAL_PROJECT_ID={{ lookup('env', 'INFISICAL_PROJECT_ID') }}
Environment=INFISICAL_TOKEN={{ lookup('env', 'INFISICAL_TOKEN') }}
Environment=INFISICAL_PATHS="/ /infrastructure/n8n /infrastructure/backup /services/vaultwarden{% if minecraft_servers is defined %}{% for server in minecraft_servers %} /services/minecraft/{{ server.id }}{% endfor %}{% endif %}"

# Start all stacks in order (traefik first, then others)
ExecStart=/bin/bash -c 'for stack in /opt/stacks/*/; do cd "$stack" && /usr/local/bin/infisical-run docker compose up -d --pull=missing; done'

# Stop all stacks
ExecStop=/bin/bash -c 'for stack in /opt/stacks/*/; do cd "$stack" && docker compose down; done'

TimeoutStartSec=300
TimeoutStopSec=120

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/stacks
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
